# ATerra/Services/avial_generals.service
[Unit]
Description=Avial Generals Adapter Service
Requires=avesterra.service
After=avesterra.service

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/AvesTerra/Executables/avial_generals
Restart=on-failure
RestartSec=2s
User=atra
Group=atra

[Install]
WantedBy=multi-user.target
