# ATerra/Services/avial_boost.service
[Unit]
Description=Avial Boost Adapter Service
Requires=avesterra.service
After=avesterra.service

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/AvesTerra/Executables/avial_boost
Restart=on-failure
RestartSec=2s
User=atra
Group=atra

[Install]
WantedBy=multi-user.target
